# work/display_setup/99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/local/bin/autorandr-systemd-user.sh"

